# backend/nas-terminal.service
[Unit]
Description=NAS Terminal Backend Agent
After=network.target docker.service

[Service]
Type=simple
WorkingDirectory=/home/ghively/remote-land/backend
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
User=ghively
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
